# etc/supervisor/conf/sshd.conf
[program:sshd]
autostart=true
autorestart=true
command=/start-sshd.sh
stdout_logfile=/var/log/sshd/out.log
stderr_logfile=/var/log/sshd/err.log
